# src/backend/src/main/resources/formfoss-backup.timer
[Unit]
Description=Run formfoss backup daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
